# scripts/systemd/sd-clone.timer
[Unit]
Description=Run sd-clone weekly

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
